# package/DEBIAN/prerm
#!/bin/bash

systemctl stop homed-zigbee.service
systemctl disable homed-zigbee.service

exit 0
